# pelagos-pfctl/assets/com.pelagos.pfctl.embedded.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label must match the binary's bundle identifier and the key in
         pelagos-mac's SMPrivilegedExecutables Info.plist. -->
    <key>Label</key>
    <string>com.pelagos.pfctl</string>

    <!-- SMJobBless installs the binary to /Library/PrivilegedHelperTools/.
         ProgramArguments must reference that path, not the source location. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/com.pelagos.pfctl</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Designated requirement for the calling binary (pelagos).
         @CALLER_DR@ is substituted at build time by pelagos-pfctl/build.rs.
         Dev default: certificate leaf[subject.CN] = "pelagos-mac Dev"
         Prod: anchor apple generic and identifier "io.pelagos.mac" and ... -->
    <key>SMAuthorizedClients</key>
    <array>
        <string>@CALLER_DR@</string>
    </array>

    <key>StandardOutPath</key>
    <string>/var/log/pelagos-pfctl.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/pelagos-pfctl.log</string>
</dict>
</plist>
